# automation/.config/systemd/user/sync-all.service
[Unit]
Description=Synchronize machine snapshot, dotfiles, notes, and wallpapers
ConditionFileIsExecutable=%h/.local/bin/sync-all

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync-all
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=30min
